# systemd/system/user@skype.service.d/env.conf
[Service]
Environment=
Environment=DISPLAY=:1
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%h/tmp/session_bus_socket

# /run/nologin stays in place for waaaay too long due to remote-fs mounts,
#  and it breaks pam login modules for systemd --user here as well, so just rm it.
ExecStartPre=-/bin/rm -f /run/nologin
