# systemd/user/dbus.socket
[Socket]
ListenStream=/home/skype/tmp/session_bus_socket
