# systemd/user/x11vnc.service
[Unit]
After=X.target

[Service]
ExecStart=/usr/bin/x11vnc -display $DISPLAY -nopw -listen localhost -forever
