# systemd/user/Xvfb-start-delay.service
[Unit]
After=Xvfb.service
Before=X.target

[Service]
Type=oneshot
ExecStart=/bin/sleep 2
